# services.d/nginx/run
#!/usr/bin/with-contenv sh
echo >&2 "starting php-fpm"

exec /usr/sbin/nginx -g "daemon off;"
