# s6/php-fpm/run
#!/bin/bash
[[ -f ./setup ]] && source ./setup

pushd /root > /dev/null
    exec su-exec root /usr/local/sbin/php-fpm 2>&1
popd
