# services.d/php-fpm/run
#!/usr/bin/with-contenv sh
echo >&2 "starting php-fpm"

exec /usr/local/sbin/php-fpm 2>&1